# linux/systemd/taskflow-backup.service
[Unit]
Description=Auto-backup TaskFlow tasks to GitHub

[Service]
Type=oneshot
ExecStart=/bin/bash %h/TaskFlow-Sync/backup.sh
